# systemd/bitcoind.service
[Unit]
Description=Bitcoin Knots (full archival node, local JSON-RPC)
Documentation=https://github.com/bitcoinknots/bitcoin/releases/tag/v29.3.knots20260508
After=network-online.target
Wants=network-online.target
# Wait for clone-safety / root expand / AP handoff in blockvase-ap before opening the chain DB.
After=blockvase-ap.service
Wants=blockvase-ap.service

[Service]
Type=simple
User=bitcoin
Group=bitcoin
# P2P gate: write args, then start with -networkactive=0|1 (see bitcoind-network-gate.sh).
# Args file is read at ExecStart (not EnvironmentFile) so ExecStartPre updates apply this boot.
ExecStartPre=+/home/blockvase/blockvase/scripts/bitcoind-network-gate.sh prepare-start
ExecStart=/bin/bash -c 'extra=$(cat /etc/bitcoin/blockvase-network-gate.args 2>/dev/null || true); [[ -n "$extra" ]] || extra="-networkactive=0"; exec /usr/local/bin/bitcoind -conf=/etc/bitcoin/bitcoin.conf -datadir=/var/lib/bitcoind $extra'
Restart=on-failure
RestartSec=30
TimeoutStartSec=infinity
LimitNOFILE=65536
PrivateTmp=true

[Install]
WantedBy=multi-user.target
